# services/sync-documents.timer
[Unit]
Description=DocumentSync

[Timer]
OnCalendar=*:0/15
Persistent=true

[Install]
WantedBy=multi-user.target
